# bin/gca-admin/upgrade-2/.config/gca-data/clients/monitor-sync.service
[Unit]
Description=Sync Monitor Service

[Service]
Type=simple
ExecStart=/usr/bin/monitor-sync.sh
Restart=always

[Install]
WantedBy=multi-user.target
