# bin/glow-monitor/monitor-sync.service
[Unit]
Description=Sync Monitor Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/monitor-sync.sh
Restart=no

[Install]
WantedBy=multi-user.target
